# docker/php-fpm/crontab/crontab
0 */6 * * * /usr/local/bin/php /var/www/bin/console  app:generate:feeds  >> var/log/cron/cron.log 2>&1



